# vlab-earlyoom/DEBIAN/postinst
#!/bin/sh

echo "Enabling and starting vlab-earlyoom.service..."
systemctl enable --now vlab-earlyoom.service || true
