# vlab-earlyoom/DEBIAN/postinst
#!/bin/sh

echo "Enabling and starting earlyoom.service..."
systemctl enable --now earlyoom.service
