# lightdm-hooks/etc/systemd/system/lightdm.service.d/hooks.conf
[Service]
Environment=LD_PRELOAD=/opt/vlab/preload/no-create-seat0.so:/opt/vlab/preload/xsessions-error-to-journald.so
